# learnCodeRecord/shell/crontab
20 17 * * * zsh /Users/fanlinbo/Documents/myself/learn/learnCodeRecord/shell/gitAction.sh
